# automated_tasks/system/hblock-unbound.timer
[Unit]
Description=Once a month update hblock and patch unbound rules

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
