# automated_tasks/archiving/archive-media-files.mypurpose.service
[Unit]
Description=archive media files for mypurpose

[Service]
Type=simple
ExecStart=-/home/jobs/scripts/by-user/root/archive_media_files.py /home/jobs/scripts/by-user/root/archive_media_files.mypurpose.yaml
User=root
Group=root

[Install]
WantedBy=multi-user.target
