# automated_tasks/archiving/archive-emails.myuser.timer
[Unit]
Description=Once a day backup myuser's emails

[Timer]
OnCalendar=*-*-* 2:20:00
Persistent=true

[Install]
WantedBy=timers.target
