# automated_tasks/archiving/pdftoocr.myuser_documents.timer
[Unit]
Description=Once a week pdftoocr myuser_documents

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
