# src/backups/android-phone-backup.myuser.timer
[Unit]
Description=Once a day backup myuser's android phone

[Timer]
OnCalendar=*-*-* 6:00:00
Persistent=true

[Install]
WantedBy=timers.target
