# system/hblock-unbound.timer
[Unit]
Description=Once every month update hblock filters for unbound

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target

